# rootfs/etc/profile.d/zz-overrides.sh
export EDITOR='micro'
export VISUAL='micro'
export SUDO_EDITOR='micro'
